# x-setup.init
#!/sbin/runscript

depend() {
	before xdm local
	need mkxf86config
}

start() {
	ebegin "Configuring OpenGL"
	/usr/sbin/x-setup
	eend $?
}
