# init.d/firmware
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

# This unpacks any firmware tarballs.

depend() {
	:
}

start() {
	ebegin "Checking for firmware to unpack"
	if [ -e /lib/firmware.tar.bz2 ]
	then
		eindent
		ebegin "Unpacking firmware"
		mkdir -p /lib/firmware
		tar xjf /lib/firmware.tar.bz2 -C /lib/firmware
		eend 0
		eoutdent
	fi
	[ -x /sbin/udevadm ] && /sbin/udevadm trigger
	eend $?
}
