# init.d/hwsetup
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

start() {
	ebegin "Starting hwsetup"
	[ -e /etc/sysconfig ] || mkdir -p /etc/sysconfig
	[ -x /usr/sbin/hwsetup ] && hwsetup
	eend $?
}
