# init.d/hwsetup
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need firmware
}

start() {
	einfo starting hwsetup
	[ -e /etc/sysconfig ] || mkdir -p /etc/sysconfig
	[ -x /usr/sbin/hwsetup ] && hwsetup
}
