# init.d/spind
#!/sbin/runscript

depend() {
	before local
}

start() {
	ebegin "Starting spin daemon"
	start-stop-daemon --start --quiet --exec /usr/sbin/spind &
	eend $?
}

stop() {
	ebegin "Stopping spind"
	start-stop-daemon --stop --quiet --exec /usr/sbin/spind
	eend $?
}
